# services/timeshift.timer
[Unit]
Description=Timeshift Backup Scheduling.

[Timer]
OnBootSec=5min
OnUnitActiveSec=1w

[Install]
WantedBy=multi-user.target
