# services/timeshift.service
[Unit]
Description=Scheduled Backup (Timeshift) 

[Service]
Type=simple
ExecStart=timeshift --create
