# services/tlp-start.service
[Unit]
Description=Start TLP automatically on boot

[Install]
WantedBy=multi-user.target

[Service]
Type=oneshot
ExecStart=tlp start && tlp usb
