# services/tlp-start.service
[Unit]
Description=Start TLP automatically on boot

[Service]
Type=simple
ExecStart=tlp start